# Startup/touchpad.conf
description "Fixes touchpad jitter and sensitivity issues"
start on startup
task
exec /home/pg/Scripts/Startup/touchpadBringup.sh
